# systemd/zephyrusctl.timer
[Unit]
Description=Periodically re-apply ryzenadj limits

[Timer]
OnActiveSec=30s
OnUnitActiveSec=30s
